# files/etc/systemd/system/lime-gsmstation-mode.service
[Unit]
Description=Lime GSM Mode switch (for scan)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/lime-gsmstation-enable-disable


